# archive/tailscale/populate_tailscale_hosts/tailscale_hosts.service
[Service]
ExecStart=/usr/local/bin/tailscale_hosts.sh

[Install]
WantedBy=default.target
